# deploy/nova-card.service
# Bot B (carte Interlace) — systemd. Copier dans /etc/systemd/system/nova-card.service
# puis : systemctl daemon-reload && systemctl enable --now nova-card
[Unit]
Description=Nova Bot B (carte Interlace) - interlace_bot
After=network-online.target docker.service nova-kyc.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/nova/interlace_bot
Environment=PYTHONPATH=/opt/nova/interlace_bot
Environment=WORKER_ID=0
ExecStartPre=/usr/bin/docker compose -p nova_il -f config/docker-compose.sandbox.yml up -d
ExecStart=/opt/nova/interlace_bot/venv/bin/uvicorn app:app --host 127.0.0.1 --port 3002
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
